# systemd/user/gpg-agent.service.d/override.conf
[Service]
# for xdg compliance
Environment=GNUPGHOME=%h/.local/share/gnupg

