# systemd/user/task-notify-due.service
[Unit]
Description=Taskwarrior due desktop notification
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/task-notify-due
