# home/.config/systemd/user/battery-threshold.service
[Unit]
Description=Apply low-battery actions

[Service]
Type=oneshot
ExecStart=%h/.local/bin/systemd/battery-threshold.sh
